# etc/systemd/system/upmpdcli.service.d/override.conf
[Service]
User=http
ExecStart=
ExecStart=/usr/bin/taskset -c 3 /usr/bin/upmpdcli -c /etc/upmpdcli.conf
