# etc/systemd/system/mpd.service.d/override.conf
[Unit]
BindsTo=mpdidle.service

[Service]
ExecStart=
ExecStart=/usr/bin/taskset -c 3 /usr/bin/mpd --no-daemon
