# etc/systemd/system/mpd.service.d/override.conf
[Unit]
BindsTo=mpdidle.service

[Service]
CPUAffinity=3
