# etc/systemd/system/mpd.service.d/override.conf
[Unit]
BindsTo=mpdidle.service

[Service]
CPUAffinity=3
ExecStart=
ExecStart=/usr/bin/mpd --systemd /srv/http/data/mpdconf/mpd.conf
