# etc/systemd/system/dnsmasq.service.d/override.conf
[Unit]
Requires=hostapd.service
After=hostapd.service
