# etc/systemd/system/hostapd.service.d/override.conf
[Unit]
Requires=startup.service
BindsTo=dnsmasq.service
