# etc/systemd/system/hostapd.service.d/override.conf
[Unit]
After=startup.service
BindsTo=dnsmasq.service
