# etc/systemd/system/nfs-server.service.d/override.conf
[Service]
ExecStart=
ExecStopPost=

ExecStart=/usr/bin/rpc.nfsd
ExecStop=/usr/bin/rpc.nfsd 0
