# etc/udev/rules.d/mouse.rules
ACTION=="add|remove", SUBSYSTEM=="input", ENV{ID_INPUT_MOUSE}=="1", RUN+="/srv/http/bash/mouse.sh"
