# etc/udev/rules.d/usbunpartitioned.rules
KERNEL=="sd[a-z]" \
ACTION=="add", \
RUN+="/srv/http/bash/settings/system.sh usbunpart"
