# etc/systemd/system/nginx.service.d/override.conf
[Unit]
After=php-fpm.service

[Service]
# allow mount
PrivateTmp=false
