# etc/systemd/system/shairport-sync.service.d/override.conf
[Service]
CPUAffinity=3
User=root
Group=root
