# etc/systemd/system/bluezdbus.service
[Unit]
Description=Bluetooth Dbus
After=bluetooth.service
Requires=bluetooth.service

[Service]
ExecStart=/srv/http/bash/bluezdbus.py
