# etc/systemd/system/cava.service
[Unit]
Description=VU level for VU LED and VU meter

[Service]
User=root
ExecStart=/srv/http/bash/vu.sh
ExecStop=/srv/http/bash/vu.sh stop
