# etc/systemd/system/lcdchar.service
[Unit]
Description=Character LCD

[Service]
ExecStart=/srv/http/bash/lcdchar.py
