# etc/systemd/system/powerbutton.service
[Unit]
Description=Power button
After=startup.service

[Service]
ExecStart=/srv/http/bash/powerbutton.sh

[Install]
WantedBy=getty.target
