# etc/systemd/system/radiowatchdog.service
[Unit]
Description=MPD webradio watchdog

[Service]
ExecStart=/srv/http/bash/radiowatchdog.sh
