# etc/systemd/system/radiofrance.service
[Unit]
Description=radiofrance metadata

[Service]
Type=simple
ExecStart=/srv/http/bash/status-radiofrance.sh
