# etc/systemd/system/bluealsa.service.d/override.conf
[Unit]
After=bluetooth.service
Requires=bluetooth.service

[Service]
CPUAffinity=3
ExecStart=
ExecStart=/usr/bin/bluealsa -p a2dp-sink -p a2dp-source
